# scripts/lego/crontab
0 0 1 * * /lego-renew.sh && nginx -s reload && echo nginx config reloaded.
